# deploy/periodic_tasks/crontab
10 * * * * /aggregations.sh >> /var/log/cron.log 2>&1
*/15 * * * * /usr/sbin/logrotate /etc/logrotate.conf
